# systemd/quantclaw.service
[Unit]
Description=QuantClaw AI Assistant
Documentation=https://github.com/openclaw/quantclaw
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=quantclaw
Group=quantclaw
WorkingDirectory=/home/quantclaw
ExecStart=/usr/local/bin/quantclaw gateway start --daemon
ExecReload=/bin/kill -HUP $MAINPID
ExecStop=/bin/kill -TERM $MAINPID
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=quantclaw

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictSUIDSGID=true
MemoryDenyWriteExecute=true
LockPersonality=true

# Resource limits
MemoryLimit=512M
CPUQuota=200%

# Environment
Environment="QUANTCLAW_LOG_LEVEL=info"
Environment="QUANTCLAW_WORKSPACE=/home/quantclaw/.quantclaw/workspace"

[Install]
WantedBy=multi-user.target